# service/jackett/run
#!/bin/bash

cd $JACKETT_DIR
exec /sbin/setuser nobody mono JackettConsole.exe --ListenPublic --Port=$JACKETT_PORT --DataFolder=$JACKETT_CONFIG
