# service/nzbget/run
#!/bin/bash

exec /sbin/setuser nobody $NZBGET_DIR/nzbget -c $NZBGET_CONFIG/nzbget.conf -s -o OutputMode=log
